# ucloud_sandboxes/systemd/ucloud-sandbox-relay.service
[Unit]
Description=UCloud Sandbox model relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ucloud
Group=ucloud
WorkingDirectory=/work/ucloud-sandboxes
Environment=UCLOUD_RELAY_HOST=0.0.0.0
Environment=UCLOUD_RELAY_PORT=8092
Environment=UCLOUD_RELAY_SANDBOX_TOKEN_FILE=/work/ucloud-sandboxes/state/relay-sandbox-token
Environment=UCLOUD_RELAY_WORKER_TOKEN_FILE=/work/ucloud-sandboxes/state/relay-worker-token
Environment=UCLOUD_RELAY_REQUEST_TIMEOUT_SECONDS=7200
Environment=UCLOUD_RELAY_WORKER_LEASE_SECONDS=600
Environment=UCLOUD_RELAY_COMPLETED_REQUEST_RETENTION_SECONDS=3600
EnvironmentFile=/etc/ucloud-sandboxes/relay.env
ExecStart=/work/ucloud-sandboxes/gateway-venv/bin/ucloud-sandboxes serve-model-relay --host ${UCLOUD_RELAY_HOST} --port ${UCLOUD_RELAY_PORT} --sandbox-bearer-token-file ${UCLOUD_RELAY_SANDBOX_TOKEN_FILE} --worker-bearer-token-file ${UCLOUD_RELAY_WORKER_TOKEN_FILE} --request-timeout-seconds ${UCLOUD_RELAY_REQUEST_TIMEOUT_SECONDS} --worker-lease-seconds ${UCLOUD_RELAY_WORKER_LEASE_SECONDS} --completed-request-retention-seconds ${UCLOUD_RELAY_COMPLETED_REQUEST_RETENTION_SECONDS}
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
